# conf.d/websockify.conf
[program:websockify]
command=/do/conf.d/websockify.sh
autorestart=true
